# var/service/sshd_chroot/run
#!/bin/sh
exec chroot /local/screencapt /usr/sbin/sshd -D

